# roles/docker/files/docker-prune.service
[Unit]
Description=Prune docker system

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -a -f
