# roles/docker/files/docker-prune.timer
[Unit]
Description=Run docker prune every day at 8:00

[Timer]
OnCalendar=*-*-* 8:00:00

[Install]
WantedBy=timers.target
